# systemd/bareos-incoming-connect@.timer.d/bareos-incoming-connect.conf
[Timer]
OnBootSec=10min
OnCalendar=minutely
